# scripts/systemd/hermes-gateway.service
[Unit]
Description=Hermes Agent Gateway (API server + messaging platforms)
Documentation=https://github.com/Devin-AXIS/hermes-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/.hermes/hermes-agent
Environment=HERMES_HOME=/root/.hermes
EnvironmentFile=-/root/.hermes/.env
# After EnvironmentFile: strip execute_code for Lumii API unless you remove this line.
# (HERMES_API_SERVER_DISABLE_CODE_EXECUTION=false in .env does NOT disable the tool — use true, or rely on this.)
Environment=HERMES_API_SERVER_DISABLE_CODE_EXECUTION=true
ExecStart=/root/.hermes/hermes-agent/venv/bin/python -m hermes_cli.main gateway run --replace
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

# Also in journal: journalctl -u hermes-gateway -f
StandardOutput=append:/root/.hermes/logs/gateway.log
StandardError=append:/root/.hermes/logs/gateway.log

[Install]
WantedBy=multi-user.target
