# scripts/systemd/hermes-gateway.service
[Unit]
Description=Hermes Agent Gateway (API server + messaging platforms)
Documentation=https://github.com/Devin-AXIS/hermes-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/.hermes/hermes-agent
Environment=HERMES_HOME=/root/.hermes
EnvironmentFile=-/root/.hermes/.env
# To disable execute_code for API /v1/runs only: set in .env
#   HERMES_API_SERVER_DISABLE_CODE_EXECUTION=true
# then: systemctl restart hermes-gateway
ExecStart=/root/.hermes/hermes-agent/venv/bin/python -m hermes_cli.main gateway run --replace
Restart=on-failure
RestartSec=5
LimitNOFILE=65535

# Also in journal: journalctl -u hermes-gateway -f
StandardOutput=append:/root/.hermes/logs/gateway.log
StandardError=append:/root/.hermes/logs/gateway.log

[Install]
WantedBy=multi-user.target
